# i3xkb.service
[Unit]
Description=i3xkb layout changer

[Service]
ExecStart=/home/kjubybot/go/bin/i3xkb
Slice=session.slice

[Install]
WantedBy=default.target
